# marytts/marytts.service
[Unit]
Description=MaryTTS Server

[Service]
Type=simple
ExecStart=/usr/bin/marytts-server
Restart=Always

[Install]
WantedBy=multi-user.target